# ansible/files/envoy.service
[Service]
ExecStart=/usr/bin/envoy -c /etc/envoy.yaml
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
